# systemd/rfid-receiver.service
[Unit]
Description=RFID Event Collector TCP Receiver
After=network.target

[Service]
Type=simple
User=rfidcollector
Group=rfidcollector
WorkingDirectory=/opt/rfid-event-collector
ExecStart=/opt/rfid-event-collector/venv/bin/python /opt/rfid-event-collector/src/main.py --config /etc/rfid-event-collector/config.json
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
